# src/modules/usb_dfu_bootloader/scripts/20-usb_bootloader.rules
SUBSYSTEM=="tty", SUBSYSTEMS=="usb", ATTRS{idVendor}=="0483", ATTRS{idProduct}=="caca", RUN+="/tmp/run_at_usb_bootloader_connect.sh", SYMLINK+="usb_bootloader"
